# detector/stinercut-detector-watcher.path
[Unit]
Description=Watch for detector start trigger file
Documentation=https://github.com/stinercut/stinercut3

[Path]
PathExists=/shared-videos/.detector-start
Unit=stinercut-detector-watcher.service

[Install]
WantedBy=multi-user.target
